# todoapp.socket
[Socket]
User=todoapp
Group=todoapp
ListenStream=80